# poruch-app/ops/poruch-health.service
[Unit]
Description=Check Poruch public readiness

[Service]
Type=oneshot
Environment=PORUCH_HEALTH_URL=https://app.munister.com.ua/readyz
ExecStart=/opt/poruch-app/ops/healthcheck.sh
